# autostart/monitors.desktop
[Desktop Entry]
Type=Application
Version=1.5
Name=monitors
Comment=Script to set up monitors
Exec=monitors
Terminal=false
