# autostart/wallpaper.desktop
[Desktop Entry]
Type=Application
Version=1.5
Name=wallpaper
Comment=Script to change wallpaper
Exec=wallpaper
Terminal=false
